# borgbackup.service
[Unit]
Description=Borg Backup

[Service]
ExecStart=/usr/local/bin/borg-backup.sh

[Install]
WantedBy=default.target

