# borgbackup.service
[Unit]
Description=Borg Backup

[Service]
Type=simple
Restart=on-failure
RestartSec=3600
ExecStart=/home/david/.local/bin/borg-backup.sh